# crontab
* * * * * /bin/bash python3 autoupdate.py