# crontab
*/5 * * * * /bin/bash python3 autoupdate.py